# x220/etc/local.d/01-zfs-import-toshiba.start
/sbin/zpool import -lR /mnt/media toshiba
